# system/etc/NetworkManager/dispatcher.d/09-timezone
#!/bin/sh
case "$2" in
    up | connectivity-change)
        timedatectl set-timezone "$(curl --fail https://ipapi.co/timezone)"
    ;;
esac
