# services/.config/systemd/user/flatpak-update.service
[Unit]
Description=Update Flatpak apps

[Service]
Type=oneshot
ExecStart=flatpak update --noninteractive
