# services/.config/systemd/user/flatpak-update.timer
[Unit]
Description=Run Flatpak update hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
