# system/etc/systemd/system/update-system-flatpaks.timer
[Unit]
Description=Daily automatic update of system Flatpaks

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
